# cronfile.cron
*/5 * * * * root python /usr/local/scripts/ast/start.py /tmp /tmp
